# deployment/systemd/celery-worker.service
[Unit]
Description=Celery Worker Service for Léguas Franzinas
After=network.target redis.service mysql.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/path/to/project/app.leguasfranzinas.pt
Environment="PATH=/path/to/venv/bin"
ExecStart=/path/to/venv/bin/celery -A my_project worker -l info --logfile=/var/log/celery/worker.log --pidfile=/var/run/celery/worker.pid --detach
ExecStop=/path/to/venv/bin/celery -A my_project control shutdown
Restart=always
RestartSec=10s

# Configurações de segurança
PrivateTmp=true
NoNewPrivileges=true

# Limites de recursos
LimitNOFILE=65536
TimeoutStopSec=600

[Install]
WantedBy=multi-user.target
