# deployment/supervisor/celery.conf
# Supervisor Configuration for Léguas Franzinas - Celery Workers

[program:leguas-celery-worker]
command=/path/to/venv/bin/celery -A my_project worker -l info
directory=/path/to/project/app.leguasfranzinas.pt
user=www-data
numprocs=1
stdout_logfile=/var/log/celery/worker.log
stderr_logfile=/var/log/celery/worker_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
priority=998

[program:leguas-celery-beat]
command=/path/to/venv/bin/celery -A my_project beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/path/to/project/app.leguasfranzinas.pt
user=www-data
numprocs=1
stdout_logfile=/var/log/celery/beat.log
stderr_logfile=/var/log/celery/beat_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
priority=999

[group:leguas-celery]
programs=leguas-celery-worker,leguas-celery-beat
priority=999
